# systemd-user/.config/systemd/user/truenas-share-mount.service
[Unit]
Description=RClone TrueNAS SMB Mount
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Lógica de fallback integrada
# Usamos rclone para montar el SMB sin ser root
ExecStart=/bin/bash -c 'if ping -c 1 -W 1 192.168.1.100 > /dev/null; then \
    REMOTE="truenas_local:share"; \
    else REMOTE="truenas_tailscale:share"; fi; \
    /usr/bin/rclone mount $REMOTE /media/stefano/truenas-share \
    --password-command "cat /home/stefano/.config/rclone/pass.txt" \
    --vfs-cache-mode writes \
    --vfs-cache-max-size 5G \
    --dir-perms 0700 \
    --file-perms 0600'

ExecStop=/usr/bin/fusermount -u /media/stefano/truenas-share
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
